# hw1/srm-empty.timer
[Unit]
Description=Run Smart Trash Can Cleanup Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
